# rootfs/etc/inittab
# inittab for uClinux
# Format:
# id:runlevels:action:process

# /etc/start must be executable script
null::sysinit:/etc/start

ttyS2::respawn:/bin/login -f root

null::wait:/bin/cat /proc/cpuinfo > /dev/tty1
